# route-proxy.service
[Unit]
Description=Route API CORS Proxy
After=network.target

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/route-proxy.py
Restart=always
User=nobody

[Install]
WantedBy=multi-user.target
